# scripts/com.example.moshi-name-sessions.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.moshi-name-sessions</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- edit this path to where you put the script -->
        <string>/Users/YOURUSER/.claude/scripts/moshi-name-sessions.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
